# docker/vsftpd/supervisord.conf
 
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
user=root

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true

[program:vsftpd]
command=/usr/sbin/vsftpd
autorestart=true
